# devel/ansible/roles/dev/files/supybot.service
[Unit]
Description=supybot
After=network-online.target
Wants=network-online.target

[Service]
Environment=GSS_USE_PROXY=yes
User=vagrant
WorkingDirectory=/vagrant
ExecStart=poetry run supybot /home/vagrant/supybot/supybot.conf

[Install]
WantedBy=multi-user.target
